# dot_config/systemd/user/monitors-changed.target
[Unit]
Description=Monitor Layout sensitive apps
