# dot_config/systemd/user/private_elecwhat.service
[Unit]
Description=Elecwhat
After=graphical-session.target
BindsTo=graphical-session.target

[Service]
ExecStartPre=/bin/sh -c 'until ping -c1 1.1.1.1 >/dev/null 2>&1; do sleep 1; done; sleep 1'
ExecStart=/usr/bin/elecwhat
Restart=on-failure

[Install]
WantedBy=graphical-session.target
