# dot_config/systemd/user/private_vesktop.service
[Unit]
Description=Vesktop (Discord client)
After=graphical-session.target
BindsTo=graphical-session.target

[Service]
ExecStartPre=/bin/sh -c 'until ping -c1 1.1.1.1 >/dev/null 2>&1; do sleep 1; done'
ExecStart=/usr/lib/vesktop/vesktop
Restart=on-failure

[Install]
WantedBy=graphical-session.target
