# dot_config/systemd/user/private_hyprpaper.service
[Unit]
Description=Hyprpaper wallpaper daemon
After=wayland-session@hyprland.desktop.target
BindsTo=wayland-session@hyprland.desktop.target

[Service]
ExecStart=/usr/bin/hyprpaper
Restart=on-failure

[Install]
WantedBy=wayland-session@hyprland.desktop.target
