# dot_config/systemd/user/private_hypridle.service
[Unit]
Description=Hypridle idle daemon
After=wayland-session@hyprland.desktop.target
BindsTo=wayland-session@hyprland.desktop.target

[Service]
ExecStart=/usr/bin/hypridle
Restart=on-failure

[Install]
WantedBy=wayland-session@hyprland.desktop.target
