# dot_config/systemd/user/private_cliphist-text.service
[Unit]
Description=Clipboard text watcher (cliphist)
After=graphical-session.target
BindsTo=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
ExecStart=/usr/bin/wl-paste --type text --watch cliphist store
Restart=on-failure

[Install]
WantedBy=graphical-session.target
