# dot_config/systemd/user/private_kanshi.service
[Unit]
Description=Kanshi dynamic display configuration
After=graphical-session.target
BindsTo=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
ExecStart=/usr/bin/kanshi
Restart=on-failure

[Install]
WantedBy=graphical-session.target
