# dot_config/systemd/user/theme-reload.target
[Unit]
Description=Theme-sensitive apps
